# packaging/whisper2text.service
[Unit]
Description=Speech to Text (whisper2text)
After=graphical-session.target
PartOf=graphical-session.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=%h/.local/share/whisperLocal/whisper2text
Restart=on-failure
RestartSec=3
Environment=QT_QPA_PLATFORM=wayland

[Install]
WantedBy=graphical-session.target
